# debian/chat/DEBIAN/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] ; then
    systemctl stop chat.service
    systemctl disable chat.service
fi
