# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ] ; then
    systemctl enable chat.service
    systemctl start chat.service
fi
